# systemd/pi-dashboard.service
[Unit]
Description=Pi Dashboard - Passenger Entertainment System
After=network.target graphical.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/usr/local/bin/pi-dashboard
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Graceful shutdown
TimeoutStopSec=30

[Install]
WantedBy=graphical.target
